# rdy-extensions/deploy/systemd-profitability/profitability-annual.service
[Unit]
Description=Trigger annual profitability review

[Service]
Type=oneshot
ExecStart=%h/scripts/profitability-trigger.sh annual_review
