# rdy-extensions/deploy/systemd-profitability/profitability-daily.service
[Unit]
Description=Trigger daily profitability brief

[Service]
Type=oneshot
ExecStart=%h/scripts/profitability-trigger.sh daily_brief
