# rdy-extensions/deploy/systemd-profitability/profitability-weekly.service
[Unit]
Description=Trigger weekly profitability digest

[Service]
Type=oneshot
ExecStart=%h/scripts/profitability-trigger.sh weekly_digest
