# rdy-extensions/deploy/systemd-profitability/profitability-monthly.service
[Unit]
Description=Trigger monthly profitability deep review

[Service]
Type=oneshot
ExecStart=%h/scripts/profitability-trigger.sh monthly_review
